# ansible/openwrt/roles/openwrt_gentoo_rootfs/templates/networkd/pxe-write-role.service.j2
[Unit]
Description=Write PXE role from DHCP lease (Option 224)
After=network-online.target
Wants=network-online.target
ConditionPathExistsGlob={{ openwrt_gentoo_networkd_lease_dir }}/*
StartLimitIntervalSec=0
StartLimitBurst=0

[Service]
Type=oneshot
ExecStart={{ openwrt_gentoo_networkd_role_script }}

[Install]
WantedBy=multi-user.target
